# apcupsd_exporter.service
[Unit]
Description=apcupsd_exporter

[Service]
User=root
ExecStart=/usr/bin/apcupsd_exporter

[Install]
WantedBy=multi-user.target
